# ssh_honeypot/on-boot-script.service
[Unit]
Description = Script run on reboot
After=default.target

[Service]
ExecStart=/usr/local/bin/scr.sh

[Install]
WantedBy=default.target
